# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Make sure logging directories exist, bail if no external storage
if findmnt /mnt/logstore; then
    mkdir -p /mnt/logstore/logs
    if ! test -d /mnt/logstore/flowdata; then
        mkdir -p /mnt/logstore/flowdata
        chown nfdump:nfdump /mnt/logstore/flowdata
        /etc/init.d/lw_nfdump restart
    fi
fi

exit 0
